# configs/launchd/com.agent-studio.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-studio.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>backend/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/opt/agent-studio</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>4936</string>
        <key>SLIDES_DIR</key>
        <string>/opt/slides</string>
    </dict>
    <key>UserName</key>
    <string>agent-studio</string>
    <key>GroupName</key>
    <string>staff</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/agent-studio/output.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/agent-studio/error.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>